# livefeed.service
[Unit]
Description=Zero-Latency Endoscopy Live Feed
Documentation=https://github.com/dreyeram/justlivefeed
After=multi-user.target systemd-logind.service
Wants=multi-user.target

# Wait for DRM/graphics to be ready
After=dev-dri-card0.device
After=dev-dri-card1.device

[Service]
Type=simple

# Give hardware time to initialize (USB capture card + DRM)
ExecStartPre=/bin/sleep 5

ExecStart=/opt/livefeed/livefeed.sh
Restart=always
RestartSec=3
TimeoutStopSec=5

# Run as root for DRM/KMS access
User=root
Group=root

# Environment
Environment=VIDEO_DEVICE=/dev/video0
Environment=WIDTH=1920
Environment=HEIGHT=1080
Environment=FRAMERATE=60
Environment=HOME=/root
Environment=XDG_RUNTIME_DIR=/run/user/0
Environment=TERM=linux

# Give access to all TTYs and DRM
SupplementaryGroups=video render input tty

# Realtime scheduling
Nice=-20
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Logging — check with: sudo journalctl -u livefeed -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=livefeed

[Install]
WantedBy=multi-user.target
