# livefeed.service
[Unit]
Description=Zero-Latency Endoscopy Live Feed
Documentation=https://github.com/prantikm07/justlivefeed
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 3
ExecStart=/opt/livefeed/livefeed.sh
Restart=always
RestartSec=2
TimeoutStopSec=5

# Run as root for direct hardware access (KMS/DRM)
User=root
Group=root

# Environment
Environment=VIDEO_DEVICE=/dev/video0
Environment=WIDTH=1920
Environment=HEIGHT=1080
Environment=FRAMERATE=60
Environment=SINK=auto
Environment=XDG_RUNTIME_DIR=/run/user/0

# Realtime scheduling
Nice=-20
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=livefeed

[Install]
WantedBy=multi-user.target
